# build/init.d/fdebug-proxy
#!/bin/sh
#
# fdebug-proxy      Startup script for the fDebug Proxy Server
#
# chkconfig: - 85 15
# description: fDebug Proxy Daemon
#
# processname: fdebug-proxy
# config: /etc/lighttpd/lighttpd.conf
# pidfile: /var/run/fdebug-proxy.pid


### BEGIN INIT INFO
# Provides: fdebug-proxy
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Short-Description: start and stop fDebug Proxy Daemon
# Description: fDebug Proxy Daemon
### END INIT INFO

# Source function library
. /etc/rc.d/init.d/functions

prog="fdebug-proxy"
fdebug="/webspace/fdebug-proxy/release/fdebug-proxy"
RETVAL=0

start() {
	echo -n $"Starting $prog: "
	daemon $fdebug
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog
	return $RETVAL
}

stop() {
	echo -n $"Stopping $prog: "
	killproc $fdebug
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog
	return $RETVAL
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	condrestart)
		if [ -f /var/lock/subsys/$prog ]; then
			stop
			start
		fi
		;;
	status)
		status $fdebug
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|condrestart|status}"
		RETVAL=1
esac

exit $RETVAL
